# 85-switch-gpi-disp.rules
ACTION=="change", KERNEL=="card1", SUBSYSTEM=="drm", RUN+="/usr/local/bin/switch_gpi_disp.sh"
